# systemd/system/eyefiserver.service.in
[Unit]
Description=EyeFi Server
Documentation=man:eyefiserver.conf(5)

[Service]
Type=simple
User=eyefi
ExecStart=@BINDIR@/eyefiserver --conf=@SYSCONFDIR@/eyefiserver.conf

[Install]
WantedBy=multi-user.target
